# lea-agent.service
[Unit]
Description=LEA SFTP Server & HI2/HI3 Receiver - Lawful Interception Agent
Documentation=https://github.com/backrouteio/lis-4g
After=network.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/root/lis-4g

ExecStart=/usr/bin/python3 /root/lis-4g/lea_sftp_server.py
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=10
StartLimitInterval=600
StartLimitBurst=5

StandardOutput=journal
StandardError=journal
SyslogIdentifier=lea-agent

LimitNOFILE=65536
LimitNPROC=32768

Environment="PYTHONUNBUFFERED=1"
Environment="TZ=Asia/Kolkata"

[Install]
WantedBy=multi-user.target
